# backend/deploy/werewolf-backend.service
[Unit]
Description=Werewolf Backend (Spring Boot)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=werewolf
Group=werewolf
WorkingDirectory=/opt/werewolf
EnvironmentFile=/etc/werewolf/env
ExecStart=/usr/bin/java -Xms256m -Xmx512m -XX:+ExitOnOutOfMemoryError -jar /opt/werewolf/werewolf-backend.jar
Restart=on-failure
RestartSec=5s
SuccessExitStatus=143
StandardOutput=journal
StandardError=journal
SyslogIdentifier=werewolf-backend

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/werewolf

[Install]
WantedBy=multi-user.target
